# jetson-agent/deploy/systemd/bebop-agent.service
[Unit]
Description=Bebop robot on-device agent (BLE provisioning + container manager)
Documentation=https://github.com/your-org/bebop
After=network-online.target bluetooth.service docker.service NetworkManager.service
Wants=network-online.target bluetooth.service docker.service NetworkManager.service
# Make sure we start after BlueZ is actually ready to accept D-Bus calls.
Requires=bluetooth.service

[Service]
Type=notify-reload
# The agent runs as root because it needs:
#   - BlueZ D-Bus access (for GATT server + advertising)
#   - NetworkManager write access (for Wi-Fi provisioning)
#   - Docker socket access (for container management)
# If you prefer tightening this, create a `bebop` user and grant
# the `bluetooth`, `netdev`, and `docker` groups instead.
User=root
Group=root

Environment=RUST_LOG=info,bebop_agent=info
Environment=BEBOP_AGENT_CONFIG=/etc/bebop/agent.toml

ExecStart=/usr/local/bin/bebop-agent
Restart=always
RestartSec=3
# Don't kill the child container supervisor tasks too aggressively.
KillMode=mixed
TimeoutStopSec=15

# Hardening — relaxed because the agent needs broad capability.
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
